# 2021/CSAW-Quals/web/no-pass-needed/supervisord.conf
[supervisord]
nodaemon=true
user=www

[program:app]
command=npm start
autorestart=true